# debian/fw-daemon.service
[Unit]
Description=Subgraph Firewall Daemon

[Service]
ExecStart=/usr/sbin/fw-daemon

[Install]
WantedBy=multi-user.target
